# other-backup-stuff/init.d/tangd-pi5
#!/sbin/openrc-run
# Tang server for Pi 5 NAS LAN auto-unlock.
# Avoids OpenRC's command_user= PAM hook by letting socat itself drop privileges.

name="Tang (Pi 5 NAS)"
description="LAN-only Tang binding daemon for Pi 5 NAS auto-unlock"

pidfile="/run/${RC_SVCNAME}.pid"
command="/usr/bin/socat"
command_args="TCP-LISTEN:${TANG_PORT:-7500},bind=${TANG_ADDR:-192.168.12.123},fork,su=tang SYSTEM:'/usr/libexec/tangd ${TANG_DIR:-/var/db/tang}'"
command_background=true

depend() {
	need net
	after firewall
}
